# docker/s6-rc.d/inspircd-is-ready/up
/etc/s6-overlay/scripts/wait-for-inspircd
