# docker/s6-rc.d/inspircd/run
#!/bin/sh
exec 2>&1
echo "Starting InspIRCd..."
exec s6-setuidgid inspircd /opt/inspircd/bin/inspircd --nofork --nolog --nopid --config /opt/inspircd/conf/inspircd.conf
